# cron/crontab
# add the following to your crontab and update the path to the fetchInbox.js file
*/5 * * * * node /email/cron/fetchInbox.js
